# launchd/io.6block.uvicorn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.6block.uvicorn</string>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/.venv/bin/uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>__HOME__/6block-data/uvicorn.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/6block-data/uvicorn.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>TZ</key>
        <string>Asia/Seoul</string>
        <key>LANG</key>
        <string>ko_KR.UTF-8</string>
        <key>LC_ALL</key>
        <string>ko_KR.UTF-8</string>
    </dict>
</dict>
</plist>
